# content/etc/init/xbmc-manual.conf
#!upstart
description "xbmc-manual"

start on (started xbmc-norun and started tty1 and started tty2)

task

script
    set +e

	splash --exit; splash --reset
	setupcon
end script
