# content/etc/init/xbmc-done.conf
start on xbmc-disabled or xbmc-failed-start or xbmc-loaded


pre-start exec renice 0 1
