# content/etc/init/xbmc-disabled.conf
start on xbian

emits xbmc-disabled

pre-start script
    grep -xq "manual" /etc/init/xbmc.override || { stop || :; exit 0; }
    initctl emit -n xbmc-disabled
end script
