# content/etc/init/xbmc-done.conf
start on xbmc-disabled or xbmc-failed-start or xbmc-loaded

emits xbmc-done

task

exec initctl emit -n xbmc-done
