# content/etc/init/xbmc-screensaver-wait.conf
start on stopped xbmc-screensaver and stopped xbmc

task

exec pkill cec-client || :
