# content/etc/init/xbmc-preload.conf
#!upstart
description "xbmc-preload"

start on filesystem and avahi

task

script
	ln -s /run/splash/splash.pid /run/sendsigs.omit.d/splash.pid
	start xbmc
end script
