# content/tmp/init/xbmc-preload.conf
#!upstart
description "xbmc-preload"

start on filesystem and avahi

task

script
	start xbmc
end script